# _docker/build/queue/supervisord.conf
[supervisord]
nodaemon=true
logfile=/proc/self/fd/1
logfile_maxbytes=0

[program:app]
directory=%(ENV_WORKDIR)s
command=%(ENV_COMMAND)s
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_WORKERS)s
autorestart=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0